# modules/observability/osquery-agent/query-packs/compliance.conf
{
  "platform": "linux,darwin",
  "queries": {
    "evidence_system_info": {
      "query": "SELECT hostname, cpu_brand, cpu_logical_cores, physical_memory FROM system_info;",
      "interval": 86400,
      "description": "compliance: host identity for evidence chain"
    },
    "evidence_os_version": {
      "query": "SELECT name, version, platform, build FROM os_version;",
      "interval": 86400
    },
    "evidence_kernel": {
      "query": "SELECT version FROM kernel_info;",
      "interval": 86400,
      "description": "compliance: kernel version evidence"
    },
    "evidence_users_with_shell": {
      "query": "SELECT uid, username, shell FROM users WHERE shell NOT IN ('/usr/sbin/nologin','/sbin/nologin','/run/current-system/sw/bin/nologin','/bin/false');",
      "interval": 86400,
      "description": "compliance: interactive accounts inventory"
    },
    "evidence_systemd_active_services_linux": {
      "query": "SELECT id, active_state FROM systemd_units WHERE type='service' AND active_state='active';",
      "interval": 86400,
      "platform": "linux"
    },
    "evidence_failed_logins_24h_linux": {
      "query": "SELECT username, time, host FROM last WHERE type IN (7, 8);",
      "interval": 86400,
      "platform": "linux",
      "description": "compliance: failed-login tracking"
    },
    "evidence_kernel_modules_linux": {
      "query": "SELECT name, used_by, status FROM kernel_modules;",
      "interval": 86400,
      "platform": "linux",
      "description": "compliance: loaded kernel modules"
    }
  }
}
