# modules/observability/osquery-agent/query-packs/inventory.conf
{
  "platform": "linux,darwin",
  "queries": {
    "system_info": {
      "query": "SELECT * FROM system_info;",
      "interval": 21600,
      "description": "host hardware + OS"
    },
    "os_version": {
      "query": "SELECT * FROM os_version;",
      "interval": 21600
    },
    "kernel_info": {
      "query": "SELECT * FROM kernel_info;",
      "interval": 21600
    },
    "kernel_modules_linux": {
      "query": "SELECT name, size FROM kernel_modules;",
      "interval": 21600,
      "platform": "linux"
    },
    "mounts": {
      "query": "SELECT device, path, type, flags FROM mounts;",
      "interval": 21600
    },
    "block_devices_linux": {
      "query": "SELECT name, parent, vendor, model, size, type FROM block_devices;",
      "interval": 21600,
      "platform": "linux"
    },
    "users": {
      "query": "SELECT uid, gid, username, shell, directory FROM users;",
      "interval": 21600
    },
    "groups": {
      "query": "SELECT gid, groupname FROM groups;",
      "interval": 21600
    },
    "systemd_units_linux": {
      "query": "SELECT id, type, active_state, sub_state FROM systemd_units WHERE type IN ('service','timer','target');",
      "interval": 21600,
      "platform": "linux"
    }
  }
}
